# supervisord.conf
[supervisord]
nodaemon=true

#[program:fahclient]
#directory=/var/lib/fahclient
#command=/docker-entrypoint.sh

[program:monit]
directory=/etc/monit
command=/etc/init.d/monit restart